# .config/profile.d/03-locale.sh
export LC_COLLATE=C.UTF-8
export LC_CTYPE=en_US.UTF-8
